# scripts/launchd/com.oncology.weekly-enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  주간 로컬 LLM 보강 LaunchAgent.
  매주 월요일 18:00 KST 에 scripts/weekly_local_enrich.sh 실행 (CI 09:00 KST 이후).
  맥이 그 시각에 자고 있으면 launchd가 다음 깨어남 시점에 한 번 따라잡아 실행한다.

  활성화:
    cp scripts/launchd/com.oncology.weekly-enrich.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.oncology.weekly-enrich.plist
  비활성화:
    launchctl unload ~/Library/LaunchAgents/com.oncology.weekly-enrich.plist
  즉시 1회 테스트:
    launchctl start com.oncology.weekly-enrich
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.oncology.weekly-enrich</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/minsung.choi/Library/CloudStorage/OneDrive-LGairesearch/VS_code/oncology_pipeline_dashboard/scripts/weekly_local_enrich.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/minsung.choi/Library/CloudStorage/OneDrive-LGairesearch/VS_code/oncology_pipeline_dashboard</string>

  <!-- launchd는 최소 환경 → python3(pyenv)·ollama·git·caffeinate 경로를 명시 -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/minsung.choi/.pyenv/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- 매주 월요일(Weekday=1) 18:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>18</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/minsung.choi/Library/CloudStorage/OneDrive-LGairesearch/VS_code/oncology_pipeline_dashboard/logs/launchd_weekly.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/minsung.choi/Library/CloudStorage/OneDrive-LGairesearch/VS_code/oncology_pipeline_dashboard/logs/launchd_weekly.err.log</string>
</dict>
</plist>
