# scripts/launchd/com.oncology.weekly-enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  주간 체크+알림 LaunchAgent (Option 2: 자동 알림 → 사용자 수동 보강).
  매주 월요일 12:00 KST 에 scripts/weekly_check_notify.sh 실행 (CI 09:00 KST + 처리시간 이후).
  신규 항목이 있으면 macOS 알림을 띄우고, 사용자가 앱을 닫고 직접 보강을 실행한다:
      bash scripts/weekly_local_enrich.sh        # MoE@44 고품질 수동 보강
  맥이 그 시각에 자고 있으면 launchd가 다음 깨어남 시점에 한 번 따라잡아 알림을 띄운다.

  활성화:
    cp scripts/launchd/com.oncology.weekly-enrich.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.oncology.weekly-enrich.plist
  비활성화:
    launchctl unload ~/Library/LaunchAgents/com.oncology.weekly-enrich.plist
  즉시 1회 테스트:
    launchctl start com.oncology.weekly-enrich
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.oncology.weekly-enrich</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/minsung.choi/Library/CloudStorage/OneDrive-LGairesearch/VS_code/oncology_pipeline_dashboard/scripts/weekly_check_notify.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/minsung.choi/Library/CloudStorage/OneDrive-LGairesearch/VS_code/oncology_pipeline_dashboard</string>

  <!-- launchd는 최소 환경 → python3(pyenv)·ollama·git·caffeinate 경로를 명시 -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/minsung.choi/.pyenv/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- 매주 월요일(Weekday=1) 12:00 (CI 09:00 + 처리시간 이후) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>12</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/minsung.choi/Library/CloudStorage/OneDrive-LGairesearch/VS_code/oncology_pipeline_dashboard/logs/launchd_weekly.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/minsung.choi/Library/CloudStorage/OneDrive-LGairesearch/VS_code/oncology_pipeline_dashboard/logs/launchd_weekly.err.log</string>
</dict>
</plist>
